# usbdload/rootfs/etc/inittab
::sysinit:/etc/init.d/rcS
::ctrlaltdel:/bin/umount -a -r
